# ctrs/wireguard/s6-init/s6-rc.d/protonvpn-killswitch/run
#!@execline@/bin/execlineb -WP

# SPDX-FileCopyrightText: 2024 Daniel Sampliner <samplinerD@gmail.com>
#
# SPDX-License-Identifier: GLWTPL

define marker /run/protonvpn-killswitch
define ips /run/protonvpn-ips
define ipCheckerUrl https://icanhazip.com
define protonServersUrl https://api.protonmail.ch/vpn/logicals

define -s curlArgs
	"--max-time 3 --retry 5 --retry-max-time 5 --fail --silent --show-error"

if { pipeline { curl $curlArgs $protonServersUrl }
	pipeline { jq --raw-output --exit-status --stream "select(.[0][4] == \"ExitIP\") | .[1]" }
	redirfd -w 1 $ips export LC_ALL C sort -u }

s6-notifyoncheck -d -s 1000 -t 500 -n 0 -c "eltest -f $marker"

loopwhilex
	if { snooze -H* -M* -S* -t $marker -T 30 }
	if { is-online }
	backtick -E ip { curl -4 $curlArgs $ipCheckerUrl }
	if { redirfd -w 1 /dev/null look $ip $ips }
	touch $marker
