# ctrs/wireguard/s6-init/s6-rc.d/wireguard/up
# SPDX-FileCopyrightText: 2024 Daniel Sampliner <samplinerD@gmail.com>
#
# SPDX-License-Identifier: GLWTPL

define conf /run/wireguard/wg0.conf

execline-umask 077

if { mkdir -p /run/wireguard }
if {
	redirfd -w 1 $conf
	foreground { cat @s6ContainerEnv@/WG_CONFIG } printf "\n" }

if { wg-quick up $conf }

loopwhilex -x 0 is-online
